# systemd/chatmail-control.service
[Unit]
Description=chatmail-control admin panel
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/chatmail-control
ExecStart=/usr/local/bin/chatmail-control serve --config /etc/chatmail-control/config.toml
Restart=always
RestartSec=5
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
