# systemd/chatmail-control.service
[Unit]
Description=Chatmail Control panel
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/chatmail-control
ExecStart=/usr/local/bin/chatmail-control serve
Restart=always
RestartSec=5
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=multi-user.target
