# systemd/chatmail-control.service
[Unit]
Description=chatmail-control admin panel
After=network.target

[Service]
Type=simple
User=chatmail-control
Group=chatmail-control
WorkingDirectory=/opt/chatmail-control
ExecStart=/usr/local/bin/chatmail-control serve --config /etc/chatmail-control/config.toml
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/chatmail-control /etc/chatmail-control

[Install]
WantedBy=multi-user.target
